# deploy/secvision.service
[Unit]
Description=SecVision FastAPI Backend
Wants=network-online.target
After=network-online.target postgresql.service

[Service]
Type=simple
User=secvision
Group=secvision
WorkingDirectory=/opt/secvision/backend
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/opt/secvision/backend/.env
ExecStart=/opt/secvision/venv/bin/python -m uvicorn main:app --host 127.0.0.1 --port 8000 --proxy-headers --forwarded-allow-ips='*'
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# 安全加固
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
